# deploy/ibgateway.service
[Unit]
Description=IB Gateway (Headless via IBC)
After=network.target

[Service]
Type=simple
User=amb
Environment="DISPLAY=:1"
# Use wrapper script to start Xvfb in background + IBC
ExecStart=/opt/amb-bot/deploy/start-ibgateway.sh
Restart=always
RestartSec=60
KillMode=mixed

[Install]
WantedBy=multi-user.target
